# example/rootfs/etc/services.d/example/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Prevent the service from restarting after it exits
# ==============================================================================

# Tell s6 not to restart the service
s6-svscanctl -t /var/run/s6/services
