# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Set environment variables from the add-on configuration
# ==============================================================================
export METER_SENSOR=$(bashio::config 'meter_sensor')
# export TOKEN=$(bashio::config 'token')
# export REGISTER_ID=$(bashio::config 'register_id')

# Start the Tibber Uploader service
python3 /usr/bin/tibber_uploader.py
