# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Set the METER_SENSOR environment variable
# ==============================================================================
export METER_SENSOR=$(bashio::config 'meter_sensor')

# Start the Tibber Uploader service
python3 /usr/bin/tibber_uploader.py
