# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Tibber Uploader service
# ==============================================================================

# Start your main application
python3 /usr/bin/tibber_uploader.py
