# example/rootfs/etc/services.d/example/run
#!/usr/bin/env bash

# Starten Sie hier Ihre Hauptanwendung
python /tibber_uploader.py
