# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Set environment variables from the add-on configuration
# ==============================================================================
export METER_SENSOR=$(bashio::config 'meter_sensor')
export EMAIL=$(bashio::config 'email')
export PASSWORD=$(bashio::config 'password')

# Log the start time and environment variables
echo "=================================================================================================="
echo "################ Starting Tibber Uploader service at $(date)" ################
echo "=================================================================================================="

# Berechnen Sie die nächste Ausführungszeit für 1:00 Uhr
next_run=$(date -d '01:00' +%s)
current_time=$(date +%s)

# Überprüfen Sie, ob die aktuelle Zeit nach 1:00 Uhr liegt
if [ "$current_time" -ge "$next_run" ]; then
    # Berechnen Sie die nächste Ausführungszeit für 1:00 Uhr am nächsten Tag
    next_run=$(date -d 'tomorrow 01:00' +%s)
fi

# Warten, bis es Zeit ist, das Skript auszuführen
while [ "$current_time" -lt "$next_run" ]; do
    sleep 60
    current_time=$(date +%s)
done

# Start the main script without redirecting stdout and stderr
python3 /usr/bin/main.py

# Log the end time
echo "=================================================================================================="
echo "################ Tibber Uploader service finished at $(date)" ################
echo "=================================================================================================="
