# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Set environment variables from the add-on configuration
# ==============================================================================
export METER_SENSOR=$(bashio::config 'meter_sensor')
export EMAIL=$(bashio::config 'email')
export PASSWORD=$(bashio::config 'password')

# Funktion, die Ihr Skript ausführt
run_my_job() {
  echo "=================================================================================================="
  echo "################ Starting Tibber Uploader service at $(date)" ################
  echo "=================================================================================================="
  python3 /usr/bin/main.py
  echo "=================================================================================================="
  echo "################ Tibber Uploader service finished at $(date)" ################
  echo "=================================================================================================="
}

# Endlosschleife
while true; do
  # Aktuelle Zeit
  current_time=$(date +%s)
  
  # Nächste Ausführungszeit: Heute um 01:00 AM oder morgen, wenn es schon vorbei ist
  next_run=$(date -d 'tomorrow 01:00' +%s)
  if [ "$current_time" -ge "$next_run" ]; then
    next_run=$(date -d '2 days 01:00' +%s)
  fi
  
  # Wartezeit bis zum nächsten Lauf in Sekunden
  wait_seconds=$((next_run - current_time))
  echo "Warte $wait_seconds Sekunden bis zum nächsten Lauf."
  
  # Warten bis zum nächsten Lauf
  sleep $wait_seconds
  
  # Führen Sie die Aufgabe aus
  run_my_job
done
